# supervisord.conf
; ============================================================
; Supervisor 进程管理配置
; 管理Flask应用和GitHub自动同步两个进程
; ============================================================

[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
logfile_maxbytes=5MB
logfile_backups=2
pidfile=/tmp/supervisord.pid

; ========== Flask应用主进程 ==========
[program:flask-app]
command=python /app/app/app.py
directory=/app/app
autostart=true
autorestart=true
startretries=3
startsecs=5
stopwaitsecs=10
stdout_logfile=/app/logs/flask-stdout.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/flask-stderr.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=3
environment=
    PYTHONUNBUFFERED="1",
    PYTHONDONTWRITEBYTECODE="1",
    FLASK_DEBUG="true",
    PORT="5000",
    DB_PATH="/app/data/ledger.db"

; ========== GitHub自动同步进程 ==========
[program:sync-worker]
command=/app/sync.sh
directory=/app
autostart=true
autorestart=true
startretries=3
startsecs=10
stopwaitsecs=15
stdout_logfile=/app/logs/sync-stdout.log
stdout_logfile_maxbytes=5MB
stdout_logfile_backups=2
stderr_logfile=/app/logs/sync-stderr.log
stderr_logfile_maxbytes=5MB
stderr_logfile_backups=2

; ========== Web界面 ==========
; Supervisor自带的Web管理界面（可选，仅内部使用）
[inet_http_server]
port=127.0.0.1:9001

[supervisorctl]
serverurl=http://127.0.0.1:9001
